# deploy/research-ops.service
[Unit]
Description=@PROJECT@ research operations control plane
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=@REPO_DIR@
EnvironmentFile=@ENV_FILE@
Environment=PYTHONUNBUFFERED=1
ExecStartPre=@PYTHON@ @REPO_DIR@/orchestrator/reconcile.py @REPO_DIR@/plan/plan.yaml
ExecStart=@PYTHON@ @REPO_DIR@/orchestrator/status_server.py
Restart=always
RestartSec=5
TimeoutStopSec=35
KillSignal=SIGTERM
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
